# build/overlay/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Faz o tty1 logar como root automaticamente (e a area de trabalho do PenaOS).
# Sem senha, porque e um sistema "live" rodando dentro da VM.
#
# FORMA PORTATIL (funciona no Arch E no Debian). NAO usar o '-o "-p -- \u"' do
# archiso aqui: no Debian o /bin/login e do shadow (nao do util-linux) e
# interpreta aquelas opcoes diferente — ele acaba PEDINDO senha e recusa
# ('Login incorrect'). Com '--autologin root' simples, o agetty chama
# 'login -f root' (ja autenticado) e entra direto nos dois sistemas.
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I $TERM
