# build/overlay/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
# Faz o tty1 logar como root automaticamente (e a area de trabalho do PenaOS).
# Sem senha, porque e um sistema "live" rodando dentro da VM.
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -- \\u' --noclear --autologin root - $TERM
